# installer/alpine/autostart
#!/sbin/openrc-run

# command_user="#WORKER#"
command_background=true
command_progress=false
command_user="root"
command="/${RC_SVCNAME}/${RC_SVCNAME}"
directory="/${RC_SVCNAME}"
error_log="/var/log/${RC_SVCNAME}.log"
output_log="/var/log/${RC_SVCNAME}.log"
pidfile="/run/${RC_SVCNAME}.pid"
